# etc/suckapy.conf
description "sucks things in python"
start on runlevel [2345]
stop on runlevel [!2345]
script
    SUCKAPY=staging
    export SUCKAPY
    exec /home/crisisnet/suckapy/venv/bin/python /home/crisisnet/suckapy/src/app.py --logto /var/log/suckapy.log
end script
respawn